# launchd/com.guardian-core.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.guardian-core</string>
    <key>ProgramArguments</key>
    <array>
        <string>{{NODE_PATH}}</string>
        <string>{{PROJECT_ROOT}}/dist/index.js</string>
    </array>
    <key>WorkingDirectory</key>
    <string>{{PROJECT_ROOT}}</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>{{HOME}}/.nix-profile/bin:{{HOME}}/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>{{HOME}}</string>
        <key>ASSISTANT_NAME</key>
        <string>Rumi</string>
        <key>ELEVENLABS_API_KEY</key>
        <string></string>
        <key>ELEVENLABS_AGENT_ID</key>
        <string></string>
        <key>ELEVENLABS_PHONE_NUMBER_ID</key>
        <string></string>
        <key>ALERT_PHONE_NUMBER</key>
        <string></string>
    </dict>
    <key>StandardOutPath</key>
    <string>{{PROJECT_ROOT}}/logs/guardian-core.log</string>
    <key>StandardErrorPath</key>
    <string>{{PROJECT_ROOT}}/logs/guardian-core.error.log</string>
</dict>
</plist>
